# systemd/user/emacs.service.d/override.conf
[Unit]
Requires=graphical-session.target
After=graphical-session.target

[Install]
WantedBy=graphical-session.target

[Service]
Environment=COLORTERM=truecolor
